# deploy/catspeed.service
[Unit]
Description=Cat Speed Trap (sensor + OLED + web dashboard)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/catspeed
# Use the project's virtualenv if you made one, else system python3.
ExecStart=/home/pi/catspeed/.venv/bin/python -m catspeed.main
Environment=CATSPEED_DB_PATH=/home/pi/catspeed/data/speed.db
# Measured 2026-07-05 (29" chord, 3 magnets @ 120°):
Environment=CATSPEED_WHEEL_DIAMETER_M=0.851
Environment=CATSPEED_MAGNETS_PER_REV=3
# Serve on port 80 so http://<hostname>/ works without :5000. Ports <1024 are
# privileged; grant just the bind capability instead of running as root.
Environment=CATSPEED_WEB_PORT=80
AmbientCapabilities=CAP_NET_BIND_SERVICE
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
